# systemd/ytsync.service.template
[Unit]
Description=YouTube Sync Audio Batch (Local GitHub Repo)
After=network-online.target

[Service]
Type=oneshot
ExecStart=REPO_PATH/scripts/yt_sync.sh
WorkingDirectory=REPO_PATH
StandardOutput=journal
StandardError=journal
