# systemd/ytsync.timer
[Unit]
Description=Timer orar pentru yt_sync

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
